# src/main/resources/config/crontab
# second minute hour day month week shell [args]
6 6 6 * * * sh sonar.sh
